# backend/deploy/tempo-backend.service
[Unit]
Description=Tempo Calendar FastAPI Backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tempo
Group=tempo
WorkingDirectory=/srv/tempo/app
EnvironmentFile=/srv/tempo/app/.env
ExecStart=/srv/tempo/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 3001
Restart=always
RestartSec=5
TimeoutStopSec=30
UMask=0077
NoNewPrivileges=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/tempo/app /srv/tempo/backups
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
