# backend/deploy/tempo-backup.timer
[Unit]
Description=Daily Tempo SQLite backup

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
