# backend/deploy/tempo-cleanup-sessions.timer
[Unit]
Description=Daily Tempo expired session cleanup

[Timer]
OnCalendar=*-*-* 03:45:00
Persistent=true

[Install]
WantedBy=timers.target
